# crontab/crontab
# launch data collection for advertising tins
0 18 * * * /home/ec2-user/data/collect/start.sh 
